# install/s6/services/30-argonaut/run
#!/usr/bin/with-contenv bash

if [ ! -f /tmp/state/30-argonaut ]; then

  echo 'Initialization Complete' >/tmp/state/30-argonaut
fi

echo '**'
echo '** Starting Argonaut Server'
exec nginx